# therasik-mcp-server/deploy/therasik-mcp.service
[Unit]
Description=TheraSIK MCP Server (Writing Suite API)
Documentation=https://github.com/jhuang7908/therasik-writing-standalone
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/therasik-mcp

# ── Environment ────────────────────────────────────────────────────────────────
# Copy .env.example → /etc/therasik-mcp/env and fill in secrets (chmod 0600)
EnvironmentFile=/etc/therasik-mcp/env

# Fallback inline (non-secret) defaults
Environment=THERASIK_MCP_PORT=8001
Environment=THERASIK_MCP_HOST=127.0.0.1
Environment=LOG_LEVEL=info

# ── Start ───────────────────────────────────────────────────────────────────────
# Uses a venv at /opt/therasik-mcp/venv created during deploy.sh
ExecStart=/opt/therasik-mcp/venv/bin/uvicorn api.main:app \
  --host ${THERASIK_MCP_HOST} \
  --port ${THERASIK_MCP_PORT} \
  --workers 2 \
  --log-level ${LOG_LEVEL}

# ── Lifecycle ───────────────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# ── Hardening (matches abenginecore pattern) ────────────────────────────────────
# KillMode=mixed keeps graceful uvicorn shutdown
KillMode=mixed
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
